# apps/observer-server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; The observer-server now owns the lifecycle of every session: it spawns
; Xvfb, x11vnc, websockify, socat and Chromium per session. Supervisord just
; keeps the API process alive.

[program:observer-server]
command=/usr/bin/node /app/apps/observer-server/dist/main.js
priority=10
autorestart=true
startsecs=3
startretries=5
stopwaitsecs=15
stdout_logfile=/tmp/observer-server.log
redirect_stderr=true
