# apps/observer-server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_UI_OBSERVER_VIEWPORT_WIDTH)sx%(ENV_UI_OBSERVER_VIEWPORT_HEIGHT)sx24 -nolisten tcp
priority=10
autorestart=true
stdout_logfile=/tmp/xvfb.log
redirect_stderr=true

[program:openbox]
command=/usr/bin/openbox --sm-disable
priority=20
autorestart=true
startretries=10
stdout_logfile=/tmp/openbox.log
redirect_stderr=true

[program:x11vnc]
; -localhost: raw VNC never leaves the container; noVNC/websockify is the only way in.
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -forever -shared -nopw -localhost -rfbport 5900 -repeat -quiet
priority=30
autorestart=true
startretries=10
stdout_logfile=/tmp/x11vnc.log
redirect_stderr=true

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 127.0.0.1:5900
priority=40
autorestart=true
stdout_logfile=/tmp/novnc.log
redirect_stderr=true

[program:cdp-proxy]
; Headed Chromium only binds CDP to 127.0.0.1; republish it on the container
; interface so the host (and only the host, via the loopback port publish)
; can attach with Playwright connectOverCDP.
command=/usr/bin/socat TCP-LISTEN:9222,fork,reuseaddr TCP:127.0.0.1:9221
priority=50
autorestart=true
stdout_logfile=/tmp/cdp-proxy.log
redirect_stderr=true

[program:observer-server]
command=/usr/bin/node /app/apps/observer-server/dist/main.js
priority=60
autorestart=true
startsecs=3
startretries=5
stopwaitsecs=15
stdout_logfile=/tmp/observer-server.log
redirect_stderr=true
